# roles/jetty/files/jetty.service
[Unit]
Description=jetty
After=
[Service]
Environment=RAILS_ENV=production
Type=simple
ExecStart=/opt/jetty/bin/jetty.sh start --ip=0.0.0.0 --allow-root
[Install]
WantedBy=multi-user.target